# packaging/systemd/ps2servers-edge.service
[Unit]
Description=PS2 Servers Edge UDPFS service
Documentation=https://github.com/NathanNeurotic/PS2-Servers/blob/main/docs/EDGE.md
After=network-online.target local-fs.target
Wants=network-online.target
# No ConditionPathIsDirectory here: the served root is configurable via
# PS2EDGE_ROOT in /etc/default/ps2servers-edge, and Edge validates --root at
# startup itself (PR #119 review: a fixed /srv/ps2 condition blocked valid
# configured mounts).

[Service]
Type=simple
User=ps2edge
Group=ps2edge
EnvironmentFile=-/etc/default/ps2servers-edge
ExecStart=/usr/local/bin/ps2servers-edge udpfs --root ${PS2EDGE_ROOT} --bind ${PS2EDGE_BIND} --port ${PS2EDGE_PORT} --data-port ${PS2EDGE_DATA_PORT} --protocol-mode ${PS2EDGE_PROTOCOL} --peer-timeout ${PS2EDGE_PEER_TIMEOUT} --read-only --log-format ${PS2EDGE_LOG_FORMAT}
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_INET AF_UNIX
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
